# dockerfiles/lgui/files/etc/service/pulseaudio/run
#!/bin/bash

exec env DISPLAY=:1 HOME=/nobody /sbin/setuser nobody pulseaudio -F /etc/xrdp/pulse/default.pa -n
