# dockerfiles/lgui/files/etc/service/autocutsel/run
#!/bin/bash

exec env DISPLAY=:1 HOME=/nobody /sbin/setuser nobody autocutsel
