# dockerfiles/danbooru/files/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -c /app/nginx.conf

[program:postgresql]
command=/usr/lib/postgresql/9.4/bin/postgres -D /var/lib/postgresql/9.4/main -c config_file=/etc/postgresql/9.4/main/postgresql.conf
user=postgres

[program:memcached]
command=/usr/bin/memcached -u danbooru

[program:rails]
command=/bin/bash -c 'PORT=9000 bundle exec unicorn'
directory=/app
user=postgres
